# examples/crontab
# run the jobwatcher script
0 8 * * *     ~/scripts/calibration_report.sh
#0 * * * *     ~/scripts/youth_bulge_report.sh
